# gost.service
[Unit]
Description=GOST SERVER
After=network.target

[Service]
Type=simple
limitDATA=102400
ExecStart=/usr/bin/gost -L henglxm:henglxm@:1113 socks5://:1113
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target
