# xmr.conf
[program:xmr]
directory=/home/ubuntu
command=/home/ubuntu/xxx
autostart=true
autorestart=true
startsecs=10
